# scripts/crontab
30 7 * * * /path/to/scripts/daily-post-major.sh
